# fab/services/templates/supervisor_elasticsearch.conf
[program:{{ project }}-{{ environment }}-elasticsearch]
command=/usr/bin/autossh -M0 -v -L 127.0.0.1:9200:127.0.0.1:9200 {{ es_endpoint }}
user={{ sudo_user }}
autostart=true
autorestart=true
stdout_logfile={{ log_dir }}/elasticsearch.log
redirect_stderr=true
stderr_logfile={{ log_dir }}/elasticsearch_error.log

