# utilities/deployment/upstart_templates/cchq_www.conf
# Upstart django script
# this script will start/stop the django server running via a management command
# In Ubuntu, symlink this script to /etc/init/
description "start and stop the cchq project"

start on started couchdb
# configuration variables.
# You'll want to change these as needed
env PROJECT_NAME=www.commcarehq.org
env DJANGO_HOME=%(code_root)s #where manage.py is
env DJANGO_PORT=8001
env DJANGO_HOST=0.0.0.0 # bind to 0.0.0.0 or other port where needed
env DJANGO_VIRTUALENV=%(virtualenv_root)s
env CP_TIMEOUT=300             # 5 minutes oughtta do it.
env PROJECT_USER=%(sudo_user)s
env SSL_CERTIFICATE=%(root)s/cchq_www.crt
env SSL_PRIVATE_KEY=%(root)s/cchq_www-key.pem
env SSL_CERTIFICATE_CHAIN=%(root)s/dimagi_cacert.crt

expect fork

pre-start script
        chdir $DJANGO_HOME
end script

script
        # Note, we're using the virtualenv's python interpreter.  Calling source/workon doesn't work here, so just call the ENV's executable instead.
        exec sudo -u $PROJECT_USER $DJANGO_VIRTUALENV/bin/python $DJANGO_HOME/manage.py runcpserver \
        host=$DJANGO_HOST \
        port=$DJANGO_PORT \
        timeout=$CP_TIMEOUT \
        ssl_certificate=$SSL_CERTIFICATE \
        ssl_private_key=$SSL_PRIVATE_KEY \
        ssl_certificate_chain=$SSL_CERTIFICATE_CHAIN &
        # the & tells upstart we're creating a daemon upstart manages PID creation for you.
end script
