# utilities/deployment/supervisor_templates/supervisor.conf
[program:cchq-formplayer]
command=java -Xmx512m -Xss1024k -classpath %(jython_home)s/jython.jar: -Dpython.home=%(jython_home)s -Dpython.executable=%(jython_home)s/bin/jython org.python.util.jython %(code_root)s/submodules/touchforms-src/touchforms/backend/xformserver.py 
user=%(sudo_user)s
autostart=true
autorestart=true
stdout_logfile=%(log_root)s/formplayer.log
redirect_stderr=true
stderr_logfile=%(log_root)s/formplayer.error.log

[program:cchq-django]
directory=%(code_root)s/
; TODO: cert stuff
command=%(virtualenv_root)s/bin/python %(code_root)s/manage.py runcpserver host=0.0.0.0 port=8001 timeout=300 
user=%(sudo_user)s
autostart=true
autorestart=true
stdout_logfile=%(log_root)s/django.log
redirect_stderr=true
stderr_logfile=%(log_root)s/django.error.log

[program:cchq-celery]
directory=%(code_root)s/
command=%(virtualenv_root)s/bin/python %(code_root)s/manage.py celeryd -v 2 -B -s celery -E -l INFO 
user=%(sudo_user)s
autostart=true
autorestart=true
stdout_logfile=%(log_root)s/celery.log
redirect_stderr=true
stderr_logfile=%(log_root)s/celery.error.log

[program:cchq-sofabed]
directory=%(code_root)s/
command=%(virtualenv_root)s/bin/python %(code_root)s/manage.py couchforms_to_django  
user=%(sudo_user)s
autostart=true
autorestart=true
stdout_logfile=%(log_root)s/sofabed.log
redirect_stderr=true
stderr_logfile=%(log_root)s/sofabed.error.log

